# deploy/quadlets/bolls-nginx.container
[Unit]
Description=Nginx reverse proxy for the production Bolls Bible deployment
After=bolls-network.service bolls-web.service bolls-edge-pod.service
Requires=bolls-network.service bolls-web.service bolls-edge-pod.service

[Container]
Image=docker.io/library/nginx:stable-alpine
ContainerName=bolls-nginx
Pod=bolls-edge.pod
CgroupsMode=no-conmon
Volume=STATIC_HOST_PATH:/home/bolls/web/static:ro,z
Volume=NGINX_MAIN_HOST_PATH:/etc/nginx/nginx.conf:ro,Z
Volume=NGINX_CONFD_HOST_PATH:/etc/nginx/conf.d:ro,Z
Volume=LETSENCRYPT_HOST_PATH:/etc/letsencrypt:ro,z
Volume=CHALLENGES_HOST_PATH:/var/www/certbot:ro,z
PodmanArgs=--cgroup-manager=cgroupfs

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=900
CPUQuota=50%
MemoryHigh=96M
MemoryMax=256M

[Install]
WantedBy=default.target
