# deploy/quadlets/bolls-db.container
[Unit]
Description=PostgreSQL for the production Bolls Bible deployment
After=bolls-back-network.service bolls-db-pod.service
Requires=bolls-back-network.service bolls-db-pod.service

[Container]
Image=docker.io/library/postgres:18-trixie
ContainerName=bolls-db
Pod=bolls-db.pod
CgroupsMode=no-conmon
Volume=POSTGRES_HOST_PATH:/var/lib/postgresql:Z
EnvironmentFile=ENV_FILE_PATH
Environment=PGDATA=/var/lib/postgresql/data/bolls
PodmanArgs=--cgroup-manager=cgroupfs --shm-size=256m
Exec=postgres -c shared_buffers=256MB -c effective_cache_size=600MB -c work_mem=16MB
HealthCmd=sh -c 'pg_isready -U "$POSTGRES_USER" -d "$POSTGRES_DB"'
HealthInterval=10s
HealthTimeout=5s
HealthRetries=10

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=900

[Install]
WantedBy=default.target
