# deploy/quadlets/bolls-db.container
[Unit]
Description=PostgreSQL for the production Bolls Bible deployment
After=bolls-back-network.service bolls-db-pod.service
Requires=bolls-back-network.service bolls-db-pod.service

[Container]
Image=docker.io/pgvector/pgvector:pg18-trixie
ContainerName=bolls-db
Pod=bolls-db.pod
CgroupsMode=no-conmon
Volume=POSTGRES_HOST_PATH:/var/lib/postgresql:Z
EnvironmentFile=ENV_FILE_PATH
Environment=PGDATA=/var/lib/postgresql/data/bolls
PodmanArgs=--cgroup-manager=cgroupfs
Exec=postgres -c shared_buffers=512MB -c effective_cache_size=2048MB -c work_mem=24MB -c maintenance_work_mem=256MB
HealthCmd=sh -c 'pg_isready -U "$POSTGRES_USER" -d "$POSTGRES_DB"'
HealthInterval=10s
HealthTimeout=5s
HealthRetries=10

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=900
CPUQuota=150%
MemoryHigh=1500M
MemoryMax=2G

[Install]
WantedBy=default.target
